# led-driver.service
[Unit]
Description=LED Display Driver Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/jim/Esp32-matrix
ExecStart=/usr/bin/python3 -m rpi_driver.main --config configs/current.json --port 8080
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
